# scripts/com.user.agentcore-light.timeline.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- AgentCore-Light timeline logger（4 源日志聚合：HOOK/API/BLE/ESP32）-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.agentcore-light.timeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__ROOT__/scripts/debug-timeline.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/agentcore-light-timeline-runner.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agentcore-light-timeline-runner.log</string>
</dict>
</plist>
