# scripts/com.user.agentcore-light.bridge.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- AgentCore-Light BLE 桥接（拉 status → 写灯）-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.agentcore-light.bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-u</string>
        <string>ble_status_bridge.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__BRIDGES_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/agentcore-light-bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agentcore-light-bridge.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
