# scripts/com.user.agentcore-light.logrotate.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- AgentCore-Light 日志轮转（每小时跑一次）-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.agentcore-light.logrotate</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__ROOT__/scripts/log-rotate.py</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/agentcore-light-logrotate.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agentcore-light-logrotate.log</string>
</dict>
</plist>
